# calendareshop/calendareshop.cron
* * * * * /www/calendar-eshop/bin/python /www/calendar-eshop/calendar-eshop/calendareshop/manage.py send_queued_mail --processes=1 >> /www/calendar-eshop/calendar-eshop/calendareshop/logs/cron_mail.log 2>&1
